# conf/supervisor/supervisord.conf
[supervisord]
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisor/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
childlogdir = /var/log/supervisor
user = root


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl = unix:///var/run/supervisor.sock


#[unix_http_server]
#file = /var/run/supervisor.sock
#chmod = 0700
#username = admin
#password = R2T9kNmEW1lvt4D8

[inet_http_server]
port = *:9001
username = admin
password = R2T9kNmEW1lvt4D8

[program:php-fpm]
command = /usr/local/sbin/php-fpm -FR
autostart = true
priority = 5
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0


[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autostart = true
priority = 10
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr	
stderr_logfile_maxbytes = 0


[include]
files = /etc/supervisor/conf.d/*.conf

